# ogame.ini
# /etc/supervisor/conf.d/ogame.ini
[program:ogame]
directory=/var/www/ogame/
command=/root/.virtualenvs/ogame/bin/python bot.py
autostart=true
autorestart=true
redirect_stderr=true
